# packages/ubuntu/debian.ntopng/postrm
#!/bin/sh -e

set -e

if [ "$1" = "purge" ] ; then
    if [ -f /bin/systemctl ] ; then
	systemctl disable ntopng.service
    else
	update-rc.d ntopng remove >/dev/null
    fi
fi

exit 0
