# packages/ubuntu/debian.ntopng/postrm
#!/bin/sh -e

set -e

if [ "$1" = "purge" ] ; then
  update-rc.d ntopng remove >/dev/null
fi

exit 0
