# packages/ubuntu/debian.ntopng/postrm
#!/bin/sh -e

set -e

case "$1" in
    purge|remove)
	update-rc.d ntopng remove >/dev/null
    ;;
esac

exit 0
