# packages/etc/systemd/system/ntopng.service
[Unit]
Description=Start/stop ntopng program
After=network.target syslog.target redis.service pf_ring.service cluster.service
Requires=pf_ring.service cluster.service

[Service]
Type=forking
RemainAfterExit=yes
ExecStart=/etc/init.d/ntopng start
ExecStop=/etc/init.d/ntopng stop
Restart=on-abnormal
RestartSec=5

[Install]
WantedBy=multi-user.target
