# packages/etc/systemd/system/ntopng.service.in
[Unit]
Description=ntopng high-speed web-based traffic monitoring and analysis tool
After=@SERVICE_AFTER@
Requires=@SERVICE_REQUIRES@

ConditionPathExists=/etc/ntopng/ntopng.conf

[Service]
Type=simple

ExecStartPre=-/bin/sh -c '/bin/sed "/-e.*$\\|-G.*\\|--daemon.*\\|--pid.*/s/^/#/" /etc/ntopng/ntopng.conf > /tmp/ntopng.conf'

ExecStart=/usr/local/bin/ntopng /tmp/ntopng.conf

ExecStopPost=-/bin/rm -rf /tmp/ntopng.conf

Restart=on-abnormal
RestartSec=5

[Install]
WantedBy=multi-user.target
