# packages/etc/init/ntopng.conf
# ntopng - 1998-13 (C) ntop.org

description "ntopng: web-based traffic monitoing"
author "Luca Deri <deri@ntop.org>"

# When to start the service
start on runlevel [2345]

# When to stop the service
stop on runlevel [016]

# Automatically restart process if crashed
respawn

# Essentially lets upstart know the process will detach itself to the background
expect fork

# Start the process
exec /etc/init.d/ntopng start

